# yocto/meta-inky-soup/recipes-inky-soup/inky-soup/files/inky-soup-server.service
[Unit]
Description=Inky Soup Web Server
After=network-online.target avahi-daemon.service
Wants=network-online.target

[Service]
Type=simple
User=inky
# Server runs from /usr/share/inky-soup where templates/static are installed.
WorkingDirectory=/usr/share/inky-soup
# Images are stored in /home/inky/inky-soup (bind-mounted from /data).
Environment=ROCKET_CLI_COLORS=off
ExecStart=/usr/bin/inky-soup-server
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
