# yocto/meta-inky-soup/recipes-support/persistent-data/files/set-hostname.service
[Unit]
Description=Set hostname from /boot/hostname.txt
DefaultDependencies=no
Before=network-pre.target systemd-hostnamed.service avahi-daemon.service
After=local-fs.target

[Service]
Type=oneshot
ExecStart=/usr/sbin/set-hostname.sh
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
